# resources/service/daemon/run
#!/bin/sh

exec chpst -u www-data /usr/bin/php /usr/local/cacti/html/cactid.php --foreground
